# APP_Source/Apps/vnctalk/input/all/etc/systemd/system/vnctalk.service
[Unit]
Description=VNCtalk

[Service]
ExecStart=/bin/bash /etc/vnctalk/vnctalk-init.sh
