# APP_Source/Apps/vncmail/input/all/etc/systemd/system/vncmail.service
[Unit]
Description=VNCmail

[Service]
ExecStart=/bin/bash /etc/vncmail/vncmail-init.sh
