# APP_Source/Tools_Drivers/brother_scanner_brscan4/input/all/etc/systemd/system/brother_scanner_brscan4.service
[Unit]
Description=Brother Scanner brscan4

[Service]
ExecStart=/bin/bash /etc/brother_scanner_brscan4/brother_scanner_brscan4-init.sh
