# APP_Source/Apps/gimp/input/all/etc/systemd/system/gimp.service
[Unit]
Description=GIMP

[Service]
ExecStart=/bin/bash /services/gimp/config/bin/gimp-init.sh
