# APP_Source/Apps/appspace/input/all/etc/systemd/system/appspace.service
[Unit]
Description= Appspace App

[Service]
ExecStart=/bin/bash /etc/appspace/appspace-init.sh