# APP_Source/Apps/kcalc/input/all/etc/systemd/system/kcalc.service
[Unit]
Description=kcalc rc file

[Service]
ExecStart=/bin/bash /etc/kcalc/kcalc-init.sh
