# APP_Source/Apps/warp/input/all/etc/systemd/system/warp.service
[Unit]
Description=Warp Terminal

[Service]
ExecStart=/bin/bash /etc/warp/warp-init.sh
