# APP_Source/Apps/beyondtrust_bomgar/input/all/etc/systemd/system/bomgar.service
[Unit]
Description=BeyondTrust Bomgar

[Service]
ExecStart=/bin/bash /etc/bomgar/bomgar-init.sh
