# APP_Source/Apps/parallels/input/all/etc/systemd/system/parallels_ras.service
[Unit]
Description=parallels_ras

[Service]
ExecStart=/bin/bash /etc/parallels_ras/parallels_ras-init.sh
