# APP_Source/Unified_Communications/ringcentral_horizon_vdi/input/all/etc/systemd/system/ringcentral_horizon_vdi.service
[Unit]
Description=RingCentral Horizon VDI Plugin

[Service]
ExecStart=/bin/bash /etc/ringcentral_horizon_vdi/ringcentral_horizon_vdi-init.sh
