# APP_Source/Apps/supertuxkart/input/all/etc/systemd/system/supertuxkart-server.service
[Unit]
Description=SuperTuxKart server

[Service]
Type=simple
DynamicUser=yes
ExecStart=/services/supertuxkart/usr/bin/supertuxkart-server.py
