# APP_Source/Apps/rclone/input/all/etc/systemd/system/rclone.service
[Unit]
Description=rclone

[Service]
ExecStart=/bin/bash /etc/rclone/rclone-init.sh
