# APP_Source/Apps/avaya_agent/input/all/etc/systemd/system/avaya_agent.service
[Unit]
Description=Avaya Agent for Desktop

[Service]
ExecStart=/bin/bash /etc/avaya_agent/avaya_agent-init.sh
