# APP_Source/Apps/intellij_idea_git/input/all/etc/systemd/system/intelij_idea_git.service
[Unit]
Description=Git User Settings Link

[Service]
ExecStart=/bin/bash /etc/intellij_idea_git/intellij_idea_git-init.sh
