# APP_Source/Apps/vnctask/input/all/etc/systemd/system/vnctask.service
[Unit]
Description=VNCtask

[Service]
ExecStart=/bin/bash /etc/vnctask/vnctask-init.sh
