# APP_Source/Apps/teamviewerqs/input/all/etc/systemd/system/teamviewerqs.service
[Unit]
Description=TeamViewer QuickSupport

[Service]
ExecStart=/services/teamviewerqs/config/bin/teamviewerqs-setup.sh