# APP_Source/Apps/mount_ios_device/input/all/etc/systemd/usbmuxd.service
[Unit]
Description=Socket daemon for the usbmux protocol used by Apple devices

[Service]
Type=oneshot
RemainAfterExit=Yes
ExecStart=/usr/sbin/usbmuxd --user usbmux --systemd