# APP_Source/Apps/vmwarerc/input/all/etc/systemd/system/vmwarerc.service
[Unit]
Description=VMware Remote Console

[Service]
ExecStart=/bin/bash /etc/vmwarerc/vmwarerc-init.sh
