# APP_Source/Apps/onedrive/input/all/etc/systemd/system/onedrive.service
[Unit]
Description=OneDrive client for Linux

[Service]
ExecStart=/bin/bash /services/onedrive/etc/onedrive/onedrive-init.sh