# APP_Source/Apps/vscode_git/input/all/etc/systemd/system/vscode_git.service
[Unit]
Description=Git User Settings Link

[Service]
ExecStart=/bin/bash /etc/vscode_git/vscode_git-init.sh
