# APP_Source/Unified_Communications/twinkle/input/all/etc/systemd/system/twinkle.service
[Unit]
Description= Twinkle

[Service]
ExecStart=/bin/bash /etc/twinkle/twinkle-init.sh