# APP_Source/Apps/ia_vision/input/all/etc/systemd/system/ia_vision.service
[Unit]
Description=Inductive Automation Vision

[Service]
ExecStart=/bin/bash /etc/ia_vision/ia_vision-init.sh
