# APP_Source/Unified_Communications/ringcentral/input/all/etc/systemd/system/ringcentral.service
[Unit]
Description=ringcentral

[Service]
ExecStart=/bin/bash /etc/ringcentral/ringcentral-init.sh
