# APP_Source/Tools_Drivers/topaz_sigplus_pro_c_object_library/input/all/etc/systemd/system/topaz.service
[Unit]
Description=Topaz SigPlus Pro C++ Object Library

[Service]
ExecStart=/bin/bash /etc/topaz/topaz-init.sh
