# APP_Source/Apps/spotify/input/all/etc/systemd/system/spotify.service
[Unit]
Description=spotify

[Service]
ExecStart=/bin/bash /etc/spotify/spotify-init.sh
