# APP_Source/Apps/keepassxc/input/all/etc/systemd/system/keepassxc.service
[Unit]
Description=KeePassXC

[Service]
ExecStart=/bin/bash /etc/keepassxc/keepassxc-init.sh
