# APP_Source/Apps/joplin/input/all/etc/systemd/system/joplin.service
[Unit]
Description=Joplin

[Service]
ExecStart=/bin/bash /etc/joplin/joplin-init.sh
