# APP_Source/Apps/joplin/input/all/etc/systemd/system/joplin.service
[Unit]
Description=Joplin

[Service]
ExecStart=/bin/bash /services/joplin/config/bin/joplin-init.sh
