# APP_Source/Apps/pcmanfm/input/all/etc/systemd/system/pcmanfm.service
[Unit]
Description=PCmanFM

[Service]
ExecStart=chmod -R 755 /services/pcmanfm
