# APP_Source/Network/f5_vpn/input/all/etc/systemd/system/f5_vpn.service
[Unit]
Description=F5 VPN

[Service]
ExecStart=/bin/bash /etc/f5_vpn/f5_vpn-init.sh
