# APP_Source/Apps/vnccontacts/input/all/etc/systemd/system/vnccontacts.service
[Unit]
Description=VNCcontacts

[Service]
ExecStart=/bin/bash /etc/vnccontacts/vnccontacts-init.sh
