# APP_Source/Apps/1password/input/all/etc/systemd/system/1password.service
[Unit]
Description=1Password

[Service]
ExecStart=/bin/bash /etc/1password/1password-init.sh
Restart=always
