# APP_Source/Apps/qztray/input/all/etc/systemd/system/qz-tray.service
[Unit]
Description=qz-tray

[Service]
ExecStart=/bin/bash /etc/qz-tray/qz-tray-init.sh
